# com.user.raycastmem.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" 
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.raycastmem</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/[YOUR_USER_NAME]/raycast_mem_monitor.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/raycastmem.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/raycastmem.err</string>
</dict>
</plist>
